# cron/cron_haf_schedule.cron.template
15 8 * * * /home/$USER/ghcontrol/scripts/Haf_Fan_Off.sh &>> $CRON_LOG 2>&1
15 19 * * * /home/$USER/ghcontrol/scripts/Haf_Fan_On.sh &>> $CRON_LOG 2>&1
