# cron/cron_haf_schedule.cron.template
15 8 * * * /home/USERNAME/ghcontrol/scripts/Haf_Fan_Off.sh &>> CRON_LOG
15 19 * * * /home/USERNAME/ghcontrol/scripts/Haf_Fan_On.sh &>> CRON_LOG
