# ghcontrol/gh_schedule.cron
@reboot /home/USERNAME/ghcontrol/log_event.sh "Started USERNAME computer"
0 * * * * /home/USERNAME/ghcontrol/get_status_8.sh

